# image_post.conf
[program:image_post]
directory=/opt/src/image_post/
command=python image_post.py
autostart=true
autorestart=true
redirect_stderr=true