# image_post.conf
[program:image_post]
directory=/opt/src/image_post/
command=python3 image_post.py
autostart=true
autorestart=false
redirect_stderr=true